# contrib/configuration/development/ramona.conf
[general]
appname={{ project_name }}
logdir=./logs/


[program:app]
directory=./api
command=gunicorn wsgi -b 127.0.0.1:9955
autorestart=yes
